# config/crontab
# TCdatalogger Automated Cron Jobs
# DO NOT EDIT MANUALLY - Generated by create_crontab.py
SHELL=/bin/zsh
PATH=/usr/local/bin:/usr/bin:/bin:/Users/shepner/src/TornCity/TCdatalogger/scripts
WORKSPACE_DIR=/Users/shepner/src/TornCity/TCdatalogger
GOOGLE_APPLICATION_CREDENTIALS=/Users/shepner/src/TornCity/TCdatalogger/config/credentials.json

*/15 * * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint server_timestamp >> /Users/shepner/src/TornCity/TCdatalogger/logs/server_timestamp.log 2>&1
*/15 * * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_40832_crimes >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_40832_crimes.log 2>&1
*/15 * * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_17991_crimes >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_17991_crimes.log 2>&1
*/15 * * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_40832_members >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_40832_members.log 2>&1
*/15 * * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_17991_members >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_17991_members.log 2>&1
0 0 * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_torn_items >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_torn_items.log 2>&1
0 */12 * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_40832_basic >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_40832_basic.log 2>&1
0 */12 * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_17991_basic >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_17991_basic.log 2>&1
0 */12 * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_40832_currency >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_40832_currency.log 2>&1
0 */12 * * * cd /Users/shepner/src/TornCity/TCdatalogger && PYTHONPATH=/Users/shepner/src/TornCity/TCdatalogger ./scripts/setup.py main --endpoint v2_faction_17991_currency >> /Users/shepner/src/TornCity/TCdatalogger/logs/v2_faction_17991_currency.log 2>&1
